# supervisord.conf
; 多源并行录制 - Supervisor 配置
; 用法: mkdir -p logs && supervisord -c supervisord.conf
; discovery_loop 每 60 秒合并新 Crypto 市场到 targets.json

[supervisord]
nodaemon=false
logfile=./logs/supervisord.log
pidfile=./logs/supervisord.pid
childlogdir=./logs

[program:discovery_loop]
command=python -u ingest/discovery_loop.py --targets targets.json --interval 60
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/discovery_loop.log
stderr_logfile=./logs/discovery_loop_err.log
environment=HTTP_PROXY="http://host.docker.internal:7897",HTTPS_PROXY="http://host.docker.internal:7897"

[program:polymarket]
command=python -u ingest/recorder.py --targets targets.json --output-dir /vault/core/data/poly/raw/polymarket
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/polymarket.log
stderr_logfile=./logs/polymarket_err.log
environment=HTTP_PROXY="http://host.docker.internal:7897",HTTPS_PROXY="http://host.docker.internal:7897"

[program:binance]
command=python -u ingest/binance_recorder.py --output-dir /vault/core/data/poly/raw/binance
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/binance.log
stderr_logfile=./logs/binance_err.log

[program:chainlink]
command=python -u ingest/chainlink_recorder.py --output-dir /vault/core/data/poly/raw/chainlink --poll-interval 0.2
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/chainlink.log
stderr_logfile=./logs/chainlink_err.log
environment=POLYGON_RPC_URL="https://polygon-rpc.com"

[program:tracker_0x8dxd]
command=python -u ingest/tracker_0x8dxd.py --output-dir /vault/core/data/poly/raw/tracker_0x8dxd --interval 30
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/tracker_0x8dxd.log
stderr_logfile=./logs/tracker_0x8dxd_err.log
environment=HTTP_PROXY="http://host.docker.internal:7897",HTTPS_PROXY="http://host.docker.internal:7897"

[program:compact]
command=python -u ingest/compact_raw_to_bronze.py --data-dir /vault/core/data/poly --interval-sec 120
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/compact.log
stderr_logfile=./logs/compact_err.log

[program:silver]
command=/bin/sh -c "while true; do python -u ingest/build_silver.py --data-dir /vault/core/data/poly --source polymarket || true; python -u ingest/build_silver.py --data-dir /vault/core/data/poly --source binance || true; python -u ingest/build_silver.py --data-dir /vault/core/data/poly --source chainlink || true; python -u ingest/build_silver.py --data-dir /vault/core/data/poly --source tracker_0x8dxd || true; sleep 300; done"
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/silver.log
stderr_logfile=./logs/silver_err.log

[program:metrics]
command=/bin/sh -c "while true; do python -u ingest/metrics_collector.py --data-dir /vault/core/data/poly --targets targets.json --output /vault/core/data/poly/checkpoint/metrics.json || true; python -u ingest/alert_check.py --metrics /vault/core/data/poly/checkpoint/metrics.json || true; sleep 60; done"
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/metrics.log
stderr_logfile=./logs/metrics_err.log

; ETL / Visual 改为手动运行：
; python -u eda/01_etl/scripts/select_top2_markets.py --data-dir /vault/core/data/poly --hours 24 --min-buys 10 --output eda/01_etl/output/top2.json --run-etl
; python -u eda/02_visual/scripts/batch_plot_retry.py --top2 eda/01_etl/output/top2.json
