# supervisord.conf
; 多源并行录制 - Supervisor 配置
; 用法: mkdir -p logs && supervisord -c supervisord.conf
; discovery_loop 每 10 分钟合并新 Crypto 市场到 targets.json

[supervisord]
nodaemon=false
logfile=./logs/supervisord.log
pidfile=./logs/supervisord.pid
childlogdir=./logs

[program:discovery_loop]
command=python -u discovery_loop.py --targets targets.json --interval 600
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/discovery_loop.log
stderr_logfile=./logs/discovery_loop_err.log
environment=HTTP_PROXY="http://127.0.0.1:7897",HTTPS_PROXY="http://127.0.0.1:7897"

[program:polymarket]
command=python -u recorder.py --targets targets.json --output-dir data/polymarket
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/polymarket.log
stderr_logfile=./logs/polymarket_err.log
environment=HTTP_PROXY="http://127.0.0.1:7897",HTTPS_PROXY="http://127.0.0.1:7897"

[program:binance]
command=python -u binance_recorder.py --output-dir data/binance
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/binance.log
stderr_logfile=./logs/binance_err.log

[program:chainlink]
command=python -u chainlink_recorder.py --output-dir data/chainlink --poll-interval 0.2
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/chainlink.log
stderr_logfile=./logs/chainlink_err.log
environment=POLYGON_RPC_URL="https://polygon-rpc.com"

[program:tracker_0x8dxd]
command=python -u tracker_0x8dxd.py --output-dir data/tracker_0x8dxd --interval 30
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=./logs/tracker_0x8dxd.log
stderr_logfile=./logs/tracker_0x8dxd_err.log
environment=HTTP_PROXY="http://127.0.0.1:7897",HTTPS_PROXY="http://127.0.0.1:7897"
